# scripts/feishu-project-agent.service
[Unit]
Description=Feishu Project Agent
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/zhaomiao/workspace-project-planning
ExecStart=/usr/bin/node src/server.js
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
